# deploy/pc-agent/com.laeddis.dingtalk-codex-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.laeddis.dingtalk-codex-agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "$HOME/dingtalk-codex-gateway" &amp;&amp; "$HOME/.local/bin/dingtalk-codex-agent" --env-file "$HOME/.config/dingtalk-codex-agent.env"</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/dingtalk-codex-agent.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/dingtalk-codex-agent.err.log</string>
</dict>
</plist>
